# crontab
SHELL=/bin/bash
DBR=$HOME/discoberry-farm

# min  hour  day  mon  dow  command
  */5  *     *    *    *    $DBR/duckdns/duck.sh > /dev/null 2>&1
