# crontab
# min  hour  day  mon  dow  command
  */5  *     *    *    *    $HOME/discoberry-farm/duckdns/duck.sh > /dev/null 2>&1
